# systemd/petrified.target
[Unit]
Description=Run Petrified Instanced Units via timers
